# scripts/auto-sync.service
[Unit]
Description=Auto-sync Final repository
After=network.target
Wants=network.target

[Service]
Type=simple
User=pc
WorkingDirectory=/home/pc/Desktop/Final
ExecStart=/home/pc/Desktop/Final/scripts/auto-sync.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
